# rpm/running-kernel.init.in
#! /bin/sh
# Copyright (c) 2004 SuSE Linux AG, Nuernberg, Germany.
# All rights reserved.
# Author: Andreas Gruenbacher <agruen@suse.de>
#
# /etc/init.d/running-kernel
#   and its symbolic link
# /usr/sbin/rcrunning-kernel
#
### BEGIN INIT INFO
# Provides:          running-kernel
# Required-Start:    $remote_fs
# Required-Stop:
# Default-Start:     2 3 5
# Default-Stop:
# Short-Description: Kernel source configuration switch
# description:  Adapt the configuration of the standard kernel sources
#		to match the currently running kernel. This is achieved
#		by creating an include file in /var/adm. The kernel sources
#		include these files to choose between the predefined
#		configurations. (This mechanism allows to build external
#		kernel modules for the standard kernels without first having
#		to configure the kernel by hand.)
### END INIT INFO

contains() {
    local w=$1
    while [ $# -ge 2 ]; do
	[ "$w" = "$2" ] && return 0
	shift
    done
    return 1
}

reconfigure() {
    local uname="$(uname -r)"
    set -- $(IFS='-' ; echo $uname)
    shift $[$#-1]

    arch=$(uname -m)
    case $arch in
    i?86)	arch=i386 ;;
    sun4u)	arch=sparc64 ;;
    arm*|sa110)	arch=arm ;;
    s390x)	arch=s390 ;;
    parisc64)	arch=parisc ;;
    esac
    # FIXME: How to handle uml?

    local flavor="$arch-$1"
    local kver="${uname%-$1}"
    # Check if this is a supported flavor. If it is not, fall back to
    # default. If there is not default configuration, fall back to any.
    if ! contains "$flavor" @FLAVORS@ ; then
	if contains $arch-default @FLAVORS@ ; then
	    flavor=$arch-default
	else
	    set -- @FLAVORS@
	    flavor="$1"
	fi
    fi

    local cfg="ARCH_FLAVOR_$(echo "$flavor" | tr a-z- A-Z_)"

    # Create the include files
    mkdir -p /var/adm
    cat > /var/adm/running-kernel.h <<-EOF
	/*
	 * Automatically generated by /etc/rc.d/running-kernel
	 *
	 * The preconfigured kernel sources include this file, and adapt their
	 * configuration to one of the following automatically:
	 *   @FLAVORS@
	 * This mechanism stops if the kernel sources are reconfigured by hand.
	 */
	#define $cfg 1
	EOF

    cat > /var/adm/running-kernel.make <<-EOF
	# Automatically generated by /etc/rc.d/running-kernel
	#
	# The preconfigured kernel sources include this file, and adapt their
	# configuration to one of the following automatically:
	#   @FLAVORS@
	# This mechanism stops if the kernel sources are reconfigured by hand.
	ARCH_FLAVOR="$flavor"
	ARCH_ALL_FLAVORS="@FLAVORS@"
	EOF

    # Adjust the arch symlink
    kdir="/usr/src/linux-$kver"
    if [ -L $kdir/include/asm -a -r $kdir/.config ] &&
       grep -q /var/adm/running-kernel.make $kdir/.config &&
       [ "$(readlink $kdir/include/asm)" != $arch ]; then
	rm -f $kdir/include/asm \
	&& ln -s asm-$arch $kdir/include/asm
    fi
}

# Status shell functions
. /etc/rc.status
# Reset status of this service
rc_reset
reconfigure
rc_exit
